# model-orchestrator/deployment/model-orchestrator.container
[Unit]
Description=Model Orchestrator - OpenAI API-compatible gateway
After=network.target

[Container]
ContainerName=model-orchestrator
Image=localhost/model-orchestrator:latest
Network=host

# Volume mounts for persistent config
Volume=%h/.config/model-orchestrator/config.yaml:/app/config.yaml:ro

# Environment variables
Environment=ROUTER_PORT=8090
Environment=ROUTER_CONFIG=/app/config.yaml

# Health check
HealthCheckCmd=/bin/sh -c "curl -f http://localhost:8090/health || exit 1"
HealthCheckInterval=30s
HealthCheckTimeout=10s
HealthCheckStartPeriod=10s
HealthCheckRetries=3

# Resource limits
Memory=2G
CPUs=2

[Install]
WantedBy=multi-user.target
